# deploy/nodesentry.service
[Unit]
Description=NodeSentry Bitcoin Core sovereignty audit
After=network-online.target bitcoind.service
Wants=network-online.target

[Service]
Type=oneshot
User=nodesentry
Group=nodesentry
ExecStart=/opt/nodesentry/bin/nodesentry audit --data-dir /var/lib/nodesentry/disk-probe --cookie-file /run/nodesentry/bitcoin-rpc-cookie --format terminal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
ReadOnlyPaths=/run/nodesentry/bitcoin-rpc-cookie /var/lib/nodesentry/disk-probe
InaccessiblePaths=-/var/lib/bitcoind/wallets -/var/lib/bitcoind/debug.log
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
